# supervisord.conf
[supervisord]
nodaemon=true

[program:sd-http]
command=prometheus-http-sd

[program:prometheus]
command=prometheus --config.file=your_config.yml